# spot_instance_terminate_notify/check-terminate
#!/bin/bash
# chkconfig: 2345 99 10
# description: spot instance terminate check

case "$1" in
    start)
        /usr/bin/terminate-notify &
        ;;
    stop) break ;;
    *) break ;;
esac
